# route53-dyndns.timer
[Unit]
Description=Dynamic DNS with Amazon Route 53

[Timer]
OnBootSec=5min
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
